# openmtc-gevent/etc/systemd/system/openmtc-backend.service
[Unit]
Description=OpenMTC Backend
Documentation=http://www.openmtc.org
Wants=ntp.service

[Service]
User=openmtc
Group=openmtc
Environment=PYTHONUNBUFFERED=true
ExecStart=/usr/local/bin/run-backend

[Install]
WantedBy=network.target
